# com.josephine.cockpit.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.josephine.cockpit.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/j0sephine/Documents/AI outputs/portal-cockpit/.venv/bin/python</string>
        <string>-m</string>
        <string>cockpit.server</string>
        <string>--port</string>
        <string>8899</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/j0sephine/Documents/AI outputs/portal-cockpit</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/j0sephine/Documents/AI outputs/portal-cockpit/logs/cockpit_server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/j0sephine/Documents/AI outputs/portal-cockpit/logs/cockpit_server.err.log</string>
</dict>
</plist>
